# data/spice-webdavd.target
[Unit]
Description=webdav daemon for Spice guests
